# deploy/halspa-runner.service
[Unit]
Description=HALSPA Test Runner Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/src/HALSPA/HALSPA-runner
ExecStart=%h/.local/bin/uv run uvicorn halspa_runner.app:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=HALSPA_RUNNER_TEST_DIR=%h/halspa-runner-duts

[Install]
WantedBy=default.target
